# etc/systemd/system/slack-logger.service
[Unit]
Description=slack-logger

[Service]
ExecStart=/opt/slack-logger/scripts/start_php

[Install]
WantedBy=multi-user.target
